# conf/ewdjs-gtm.conf
description "EWD.js server"
author      "Mike Clayton on behalf of http://mgateway.com"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 10 5

env HOME=/home/ubuntu
env NODEFOLDER=/home/ubuntu/ewdjs
env PATH=/home/ubuntu/.nvm/v0.10.28/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/usr/games:/usr/local/games

script

    chdir $NODEFOLDER

    exec bash -c 'source /home/ubuntu/.nvm/nvm.sh && exec node ewdStart-gtm.js gtm-config' >> /var/log/ewdjs.log 2>&1

end script

post-start script
   # Optionally put a script here that will notifiy you node has (re)started
end script
