# extras/conf/ewdjs-dewdrop.conf
description "EWD.js Upstart script for dEWDrop server"
author      "Rob Tweed"

start on started mountall
stop on shutdown

# Automatically Respawn:
respawn
respawn limit 10 5

# Change the value of USER below to the appropriate user name for your dEWDrop Server

env USER="vista"

script
    exec su -l $USER -c 'source ~/.nvm/nvm.sh && cd ~/ewdjs && node ewdStart-gtm.js dewdrop-config' >> /var/log/ewdjs.log 2>&1
end script


